# systemd/dnstt.service
[Unit]
Description=dnstt DNS Tunnel Server
After=network.target

[Service]
Type=simple
ExecStart=/opt/dnstt/dnstt-server \
  -udp :53 \
  -privkey-file /opt/dnstt/server.key \
  t.domain.com \
  127.0.0.1:1080

Restart=always
RestartSec=5

User=root
WorkingDirectory=/opt/dnstt

[Install]
WantedBy=multi-user.target